# sentienteye.service
[Unit]
Description=Sentient Eye
AFter=network.target graphical.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/projects/SentientEye
ExecStart=/home/pi/projects/SentientEyeEnv/bin/python sentienteye_main.py
Restart=always
RestartSec=5

Environment=DISPLAY=:0

[Install]
WantedBy=graphical.target